# systemd/printfarm-prusalink.service
[Unit]
Description=O.D.I.N. PrusaLink Monitor
After=network.target printfarm-backend.service
Wants=printfarm-backend.service

[Service]
Type=simple
WorkingDirectory=/opt/printfarm-scheduler/backend
ExecStart=/opt/printfarm-scheduler/backend/venv/bin/python3 prusalink_monitor.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
